# scripts/toss-autodeploy.timer
[Unit]
Description=Check GitHub for Toss updates every minute

[Timer]
OnBootSec=30
OnCalendar=*:0/1
AccuracySec=10
Persistent=true
Unit=toss-autodeploy.service

[Install]
WantedBy=timers.target
